# packages/gitp/gitp-0.2.5-py3-none-any.whl/etc/profile.d/git.sh
# Git system config
export GIT_CONFIG_SYSTEM='/etc/gitconfig'

# Git config to store repositories and hidden data
export GIT_CONFIG_GLOBAL='/etc/.gitconfig'
